# container_restart.timer
[Unit]
Description=Run Docker Log Monitor every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target